# build/rpm/etc/systemd/system/kbnd.service
[Unit]
Description= kbnd Service

[Service]
Type=forking
ExecStart=/etc/init.d/kbnd start
ExecStop=/etc/init.d/kbnd stop

[Install]
WantedBy=multi-user.target