# build/rpm/etc/systemd/system/kscnd.service
[Unit]
Description= kscnd Service

[Service]
Type=forking
ExecStart=/etc/init.d/kscnd start
ExecStop=/etc/init.d/kscnd stop

[Install]
WantedBy=multi-user.target